# scripts/claude-mem-daemon.service.template
[Unit]
Description=Memlayer - Claude Code Memory Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{DAEMON_PATH}}
Restart=always
RestartSec=10
EnvironmentFile=%h/.config/memlayer/env
Environment=MEMLAYER_WATCH_PATH=%h/.claude/projects
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
